# config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
user=root
logfile=/config/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid
loglevel=info

; ── 1. Xorg (real Intel GPU via modesetting + DRI3, headless HDMI dummy) ─────
[program:xorg]
command=Xorg :1 vt7 -config /etc/retroshine/xorg.conf -nolisten tcp -noreset -logfile /config/logs/xorg.log
priority=10
autorestart=true
startsecs=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; ── 2. PulseAudio (virtual null sink — captured by Sunshine) ─────────────────
; --daemonize=no keeps PulseAudio in the foreground so supervisord tracks it
; correctly. Without it PulseAudio forks, supervisord considers the parent
; exited and restarts — the child still holds the socket, causing every
; subsequent start to fail with "Address already in use".
; The rm -f guards against a stale socket left by SIGKILL on a previous stop.
[program:pulseaudio]
command=bash -c "rm -f /tmp/runtime/pulse/native && exec pulseaudio --daemonize=no --disallow-exit --disallow-module-loading=false --log-target=file:/config/logs/pulseaudio.log --exit-idle-time=-1 --file=/etc/retroshine/pulse.pa"
environment=HOME="/root",PULSE_RUNTIME_PATH="/tmp/runtime/pulse",XDG_RUNTIME_DIR="/tmp/runtime"
priority=20
autorestart=true
startsecs=3
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; ── 3. Sunshine (streaming server) ───────────────────────────────────────────
[program:sunshine]
command=/usr/bin/sunshine /config/sunshine.conf
priority=40
autorestart=true
startsecs=5
environment=DISPLAY=":1",PULSE_SERVER="unix:/tmp/runtime/pulse/native",LIBVA_DRIVER_NAME="iHD",HOME="/root",XDG_RUNTIME_DIR="/tmp/runtime"
stdout_logfile=/config/logs/sunshine.log
stderr_logfile=/config/logs/sunshine.log
stopasgroup=true
killasgroup=true
